# charon/etc/systemd/system/mnt.automount
[Unit]
Description=Automount orpheus:

[Automount]
Where=/mnt
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
